# onkyori_listener2.service
#
# Author: jun10000 (https://github.com/jun10000)
#

# Please copy to /usr/lib/systemd/system/

[Unit]
Description=Onkyo RI Listener 2
After=mariadb.service

[Service]
ExecStart=/var/project/onkyori_listener2/index.py
Nice=-20
Restart=always
Type=forking
PIDFile=/var/run/OnkyoRIListener2.pid

[Install]
WantedBy=multi-user.target
